# infrastructure/ansible/roles/app/templates/portfolio.service.j2
# ============================================================================
# PORTFOLIO SYSTEMD SERVICE
# ============================================================================
# Generated by Ansible - Do not edit manually
#
# This service ensures the portfolio container restarts on system reboot.
# Note: Docker's --restart=unless-stopped should handle this, but this
# provides additional control and logging.
# ============================================================================

[Unit]
Description=Portfolio Application Container
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=root

# Start command - ensures container is running
ExecStart=/bin/bash -c 'docker start {{ app_name | default("portfolio") }} 2>/dev/null || echo "Container not found, run deploy.sh first"'

# Stop command
ExecStop=/usr/bin/docker stop {{ app_name | default("portfolio") }}

# Restart policy
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
